# overlay/etc/rc.local
#!/bin/bash

function build_conf {
	set -e
	MY_IP_ADDR=$(/usr/local/bin/oc-metadata --cached PUBLIC_IP_ADDRESS)
	MY_IP6_ADDR=$(/usr/local/bin/oc-metadata --cached IPV6_ADDRESS)

	cat <<EOF > /root/client.ovpn
client
nobind
dev tun
redirect-gateway def1
<key>
`cat /etc/openvpn/key.pem`
</key>
<cert>
`cat /etc/openvpn/cert.pem`
</cert>
<ca>
`cat /etc/openvpn/cert.pem`
</ca>
<dh>
`cat /etc/openvpn/dh.pem`
</dh>
`[ "$MY_IP6_ADDR" != "" ] && echo -e "<connection>\nremote $MY_IP6_ADDR udp6\n</connection>"`
<connection>
remote $MY_IP_ADDR 1194 udp
</connection>
`[ "$MY_IP6_ADDR" != "" ] && echo -e "<connection>\nremote $MY_IP6_ADDR tcp6-client\n</connection>"`
<connection>
remote $MY_IP_ADDR 443 tcp-client
</connection>
EOF
	PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

	BASEDIR=/var/www/openvpn
	UUID=$(uuid)

	rm -rf $BASEDIR && mkdir -p $BASEDIR/$UUID && touch $BASEDIR/index.html
	cp /root/client.ovpn $BASEDIR/$UUID/$(oc-metadata --cached HOSTNAME).ovpn

	scw-sync-kernel-modules
	depmod -a
}

# OpenVPN Setup
[ -d /dev/net ] || \
    mkdir -p /dev/net
[ -c /dev/net/tun ] || \
    mknod /dev/net/tun c 10 200
[ -f /etc/openvpn/dh.pem ] || \
    openssl dhparam -out /etc/openvpn/dh.pem 2048
[ -f /etc/openvpn/key.pem ] || \
    openssl genrsa -out /etc/openvpn/key.pem 2048
chmod 600 /etc/openvpn/key.pem
[ -f /etc/openvpn/csr.pem ] || \
    openssl req -new -key /etc/openvpn/key.pem -out /etc/openvpn/csr.pem -subj /CN=OpenVPN/
[ -f /etc/openvpn/cert.pem ] || \
    openssl x509 -req -in /etc/openvpn/csr.pem -out /etc/openvpn/cert.pem -signkey /etc/openvpn/key.pem -days 24855
[ -f /root/client.ovpn ] || \
	build_conf && \
	[ "$(oc-metadata --cached IPV6_ADDRESS)" != "" ] && cat /etc/openvpn/udp1194-ipv6.conf >>/etc/openvpn/udp1194.conf && \
	[ "$(oc-metadata --cached IPV6_ADDRESS)" != "" ] && cat /etc/openvpn/tcp443-ipv6.conf >>/etc/openvpn/tcp443.conf && \
	sleep 10 && \
    systemctl enable openvpn@tcp443 && \
    systemctl enable openvpn@udp1194 && \
	systemctl start openvpn@tcp443 && \
    systemctl start openvpn@udp1194 && \
	sed -i "s/\$ID/$(oc-metadata --cached ID)/g" /etc/caddy/Caddyfile &&\
	systemctl enable caddy && \
	systemctl start caddy

#v6 Setup
[ "$(oc-metadata --cached IPV6_ADDRESS)" != "" ] && /usr/local/bin/nat6 &

exit 0